# docs/deploy/demo-publica/rhsaas-demo-pool-maintenance.timer
[Unit]
Description=RH SaaS - manutencao periodica da pool demo

[Timer]
OnBootSec=2min
OnUnitActiveSec=2min
RandomizedDelaySec=10s
AccuracySec=15s
Persistent=true
Unit=rhsaas-demo-pool-maintenance.service

[Install]
WantedBy=timers.target
